# data/systemd/icd.service.in
[Unit]
Description=Image Capture Microservice
AssertPathExists=@sysconfdir@

[Service]
#User=pi
#Group=pi
ExecStart=@bindir@/icd --config @sysconfdir@/icd.conf

[Install]
WantedBy=default.target
